# build/plist/com.flyenv.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.flyenv.helper</string>

    <key>StandardErrorPath</key>
    <string>/var/log/com.flyenv.helper.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/com.flyenv.helper.log</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Application Support/FlyEnv/Helper/flyenv-helper</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
