# build/plist/com.flyenv.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.flyenv.helper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/FlyEnv.app/Contents/Resources/helper/flyenv-helper</string>
        <string>/Applications/FlyEnv.app/Contents/Resources/helper/helper.js</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
